# iso/profiles/common/airootfs/etc/systemd/system/foundationos-live-prepare.service
[Unit]
Description=Prepare FoundationOS live environment
After=local-fs.target systemd-user-sessions.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/foundationos-live-prepare
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
