# iso/profiles/common/airootfs/etc/profile.d/foundationos-live.sh
export FOUNDATIONOS_ROOT=/opt/foundationos-client
export FOUNDATIONOS_PATH=/opt/foundationos-client
export OMARCHY_PATH="${OMARCHY_PATH:-$HOME/.local/share/omarchy}"
export PATH=/opt/foundationos-client/bin:$PATH
