# config/systemd/user/foundation-desktop-agent.service
[Unit]
Description=Foundation Desktop Agent
After=graphical-session.target network-online.target
Wants=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
EnvironmentFile=-%h/.config/foundationos/agent.env
ExecStart=/usr/bin/env bash -lc '${FOUNDATION_AGENT_BINARY:-foundation-desktop-agent}'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
